# smf/manifests/cns-updater.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2014, Joyent, Inc.
-->

<service_bundle type="manifest" name="cns-updater">
  <service name="triton/application/cns-updater" type="service" version="1">

    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/network/physical" />
    </dependency>

    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>

    <dependency name="mdata"
                grouping="require_all"
                restart_on="none"
                type="service">
        <service_fmri value="svc:/smartdc/mdata:execute" />
    </dependency>

    <dependency name="config-agent"
                grouping="optional_all"
                restart_on="none"
                type="service">
        <service_fmri value="svc:/smartdc/application/config-agent" />
    </dependency>

    <exec_method type="method"
      name="start"
      exec="node --abort_on_uncaught_exception updater.js &amp;"
      timeout_seconds="30">
      <method_context working_directory="/opt/triton/cns">
        <method_credential user="cns"
          group="nobody"
          privileges="basic"/>
        <method_environment>
          <envvar name="PATH"
            value="/opt/triton/cns/build/node/bin:/opt/local/bin:/usr/bin:/usr/sbin:/bin"/>
          <envvar name="LD_PRELOAD_32"
            value="/usr/lib/extendedFILE.so.1" />
        </method_environment>
      </method_context>
    </exec_method>

    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="30" />

    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='contract' />
    </property_group>

    <instance name="default" enabled="true" />

    <stability value='Unstable' />

    <template>
      <common_name><loctext xml:lang="C">Triton CNS updater</loctext></common_name>
    </template>

  </service>
</service_bundle>
